# beersomme/beersomme.conf
[program:myserver]
command=gunicorn run:app -w 1 -b 0.0.0.0:80 --error-logfile /home/ubuntu/errorlog.log

[supervisord]
logfile=/home/ubuntu/beersomme.log
loglevel=debug
user=root

